# infra/hetzner/systemd/eight-bridge-eightgent.service
[Unit]
Description=Eight Telegram bridge (@eightgentcodebot)
After=eight-daemon.service network-online.target
Wants=eight-daemon.service network-online.target

[Service]
Type=simple
Environment=BOT_NAME=eightgent
Environment=EIGHT_TG_SESSIONS=/var/lib/8gent/telegram-sessions-eightgent.json
EnvironmentFile=/etc/8gent/daemon.env
EnvironmentFile=/etc/8gent/bots.env
WorkingDirectory=/opt/8gent
ExecStart=/usr/local/bin/bun run packages/daemon/telegram-bridge.ts
Restart=always
RestartSec=5
StandardOutput=append:/var/log/8gent/bridge-eightgent.log
StandardError=append:/var/log/8gent/bridge-eightgent.log

[Install]
WantedBy=multi-user.target
